# config/com.gate.server.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gate.server</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Replace with your Gate installation path -->
        <string>/path/to/gate/.venv/bin/gate</string>
        <string>up</string>
        <string>--headless</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!--
      Runtime data (logs, state, socket) lives under platformdirs'
      user_data_dir: ~/Library/Application Support/gate on macOS.
      Point stdout/stderr there too so LaunchAgent output joins them.
      Replace USERNAME with your macOS username.
    -->
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Application Support/gate/logs/server-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Application Support/gate/logs/server-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/path/to/gate</string>
</dict>
</plist>
